# tools/nebula-maintainer/nebula-hourly.service
[Unit]
Description=Nebula hourly research + upgrade generation

[Service]
Type=oneshot
WorkingDirectory=/Users/josh/.openclaw/workspace/projects/nebula
ExecStart=/usr/bin/env bash -lc 'npm run nebula:hourly'
